# supervisor/magellan.conf
[program:magellan]
directory=/opt/magellan/src
command=python magellan.py
autorestart=false
startsecs=0

[program:magellan_web]
directory=/opt/magellan/src
command=python magellan_web.py 80
