# supervisor/magellan.conf
[program:magellan]
directory=/opt/magellan/src
command=python magellan.py

[program:magellan_web]
directory=/opt/magellan/src
command=python magellan_web.py 80
