# supervisor/register_in_service_discovery.conf
[program:register_in_service_discovery]
directory=/opt/microservice/src
command=python register_in_service_discovery.py 80 --single-active-instance
startsecs=0
priority=10
